# dhcp-wlan0.service
[Unit]
Description=DHCP for wlan0
After=wifi.service

[Service]
ExecStart=/sbin/dhclient wlan0

[Install]
WantedBy=multi-user.target
